# .docker/supervisor/webdav.conf
[program:webdav]
process_name=webdav
command=start-webdav.sh
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
autostart=%(ENV_JUPYTER_ENABLE_WEBDAV)s
