# supervisor/gateway.conf
[program:gateway]
process_name=gateway
command=start-gateway.sh
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
autostart=true
autorestart=true
exitcodes=0
